# scripts/systemd/hamsterfi-web.service
[Unit]
Description=Hamster‑Fi Web UI
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/hamster-fi
ExecStart=/opt/hamster-fi/.venv/bin/uvicorn hamsterfi.main:app --host 0.0.0.0 --port 8080
Restart=always

[Install]
WantedBy=multi-user.target
