# scripts/systemd/hamsterfi-apply.service
[Unit]
Description=Apply Hamster‑Fi config at boot
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/hamster-fi
ExecStart=/opt/hamster-fi/.venv/bin/python -c "from hamsterfi.core.config import load_config; from hamsterfi.system.apply import apply; apply(load_config())"

[Install]
WantedBy=multi-user.target
